# sdks/go/node/containerruntime/containerd/examples/containerd.service.d/http-proxy.conf
# /etc/systemd/system/containerd.service.d/http-proxy.conf
#
# containerd performs the registry pulls, so it (not just the build shell) must
# know to bypass the mandatory forward proxy for the ECR/S3 VPC endpoints, which
# are reachable directly. Without NO_PROXY here, pulls to the PTC fail.
#
# After installing: `systemctl daemon-reload && systemctl restart containerd`.
[Service]
Environment="HTTP_PROXY=http://forward-proxy:3128"
Environment="HTTPS_PROXY=http://forward-proxy:3128"
Environment="NO_PROXY=*.dkr.ecr.us-west-2.amazonaws.com,api.ecr.us-west-2.amazonaws.com,*.s3.us-west-2.amazonaws.com,localhost,127.0.0.1"
